# board/coreboot/rootfs-additions/etc/network/if-pre-up.d/00macchanger
#!/bin/sh

if [ "$IFACE" != "lo" ]; then
	/usr/bin/macchanger -A "$IFACE"
fi
