# board/coreboot/rootfs-additions/etc/profile.d/aliases.sh
alias ll='ls -lhF'
alias cls='clear;reset'
